# systemd/dacroq-data-api.service
[Unit]
Description=Dacroq Data API - Database and Data Management
After=network.target
Wants=network.target

[Service]
Type=simple
User=bdatsko
Group=bdatsko
WorkingDirectory=/var/www/dacroq/data/api
Environment=PATH=/home/bdatsko/.local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/var/www/dacroq/data/api:/home/bdatsko/.local/lib/python3.9/site-packages
ExecStart=/usr/bin/python3 /var/www/dacroq/data/api/app.py --port 8001
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/dacroq/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target 